# com.eliran.linkedinjobs.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS LaunchAgent — runs run.sh on an interval.

  Install:
    cp com.eliran.linkedinjobs.plist ~/Library/LaunchAgents/
    launchctl load   ~/Library/LaunchAgents/com.eliran.linkedinjobs.plist
    launchctl start  com.eliran.linkedinjobs          # kick it off now

  Check status:
    launchctl list | grep linkedin

  Stop:
    launchctl unload ~/Library/LaunchAgents/com.eliran.linkedinjobs.plist

  Current schedule: every 12 hours (43200s).
  To change to twice daily at specific times, swap StartInterval for
  StartCalendarInterval (see comment below).

  IMPORTANT: search.py launches a *visible* Chromium window (headless=False is
  required for LinkedIn). That means you'll see a browser pop up during
  scheduled runs. If you want it silent, change headless=False to True in
  search.py — LinkedIn may detect headless and block the scrape; test first.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.eliran.linkedinjobs</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/eliranei/linkedin-jobs/run.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/eliranei/linkedin-jobs</string>

  <!-- Every 12 hours. -->
  <key>StartInterval</key>
  <integer>43200</integer>

  <!-- ALTERNATIVE: fixed times. Uncomment and delete StartInterval above.
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>0</integer></dict>
  </array>
  -->

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/eliranei/linkedin-jobs/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/eliranei/linkedin-jobs/launchd.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <!-- Backend mode for unattended cron runs. "guest" = browserless HTTP
         (recommended). "loggedin" = Playwright + saved session (only useful
         if you want to scrape from a logged-in dummy account). -->
    <key>LINKEDINJOBS_MODE</key>
    <string>guest</string>
  </dict>
</dict>
</plist>
